# iems5722_app_a4.conf
[program:assignment_4]
command = /root/IEMS5722-A4/venv/bin/python3 assignment_4.py
directory = /root/IEMS5722-A4
user = root
autostart = true
autorestart = true
stdout_logfile = /root/IEMS5722-A4/app.log
redirect_stderr = true
